# service/openresty/run
#!/bin/sh
openresty -p /etc/nginx/ -c /etc/nginx/nginx.conf