# cron.d/bot_blocker
0 3,15 * * * root /usr/sbin/update-ngxblocker && openresty -s reload
