# cron.d/spamhaus_blocker
0 */12 * * * root /opt/scripts/spamhaus.sh && openresty -s reload
